# hub/11-shelly-monitor/plant-shelly-monitor.service
[Unit]
Description=Plant Autonomy Testbed - sample Shelly uptime/RSSI
After=plant-listener.service

[Service]
Type=oneshot
User=basilpi
EnvironmentFile=/etc/plant-hub/credentials
WorkingDirectory=/home/basilpi/plant-hub
ExecStart=/usr/bin/python3 /home/basilpi/plant-hub/shelly_monitor.py
